# deploy/meetingnotes-fetch.service
[Unit]
Description=MeetingNotes — fetch raw pad content from Riseup
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/opt/meetingnotes
# Tuesday meeting → fetch Wednesday morning; "yesterday" gives Tuesday's date
ExecStart=/bin/bash -c '\
  DATE=$(date -d yesterday +%%Y_%%m_%%d) && \
  /opt/meetingnotes/venv/bin/python3 -m pipeline.run \
    --date "$DATE" --fetch-only'
EnvironmentFile=/opt/meetingnotes/.env
StandardOutput=journal
StandardError=journal
